# excercises/Ex03_ELT_pipeline_docker/runmain.crontab
*/1 * * * * python /workdir/scripts/main.py
